# core/installation/artemis.service
[Unit]
Description=artemis startup 

[Service]
ExecStart=/bin/artemisStart

[Install]
WantedBy=multi-user.target