# core/installation/pwnhyve.service
[Unit]
Description=pwnhyve startup
Before=usb-gadget.target

[Service]
ExecStart=/bin/pwnhyveStart

[Install]
RequiredBy=usb-gadget.target