# core/installation/pwnhyve.service
[Unit]
Description=pwnhyve startup 

[Service]
ExecStart=/bin/pwnhyveStart

[Install]
WantedBy=multi-user.target